# service/.config/systemd/user/brightness.service
[Unit]
Description=Notify when the battery is low and the computer is discharching.

[Service]
ExecStart=brightnessctl s 1%
