# service/.config/systemd/user/dayshift.service
[Unit]
Description=Switching day

[Service]
ExecStart=~/.local/share/bin/dayshift

