# service/.config/systemd/user/dayshift.service
[Unit]
Description=Switching day

[Service]
ExecStart=sh -c ~/.local/bin/dayshift

[Install]
WantedBy=graphical.target
